# files/cron.d/puppet
0 */3 * * * root /usr/bin/puppet apply /etc/setup-files/manifests/default.pp
